# parentalparrot.service
[Unit]
Description=ParentalParrot - Device usage monitor and limiter
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/ParentalParrot -db /var/lib/parentalparrot/parentalparrot.db -seed /etc/parentalparrot/config.example.yaml
Restart=on-failure
RestartSec=10
WorkingDirectory=/var/lib/parentalparrot

[Install]
WantedBy=multi-user.target
